# LaunchAgents/com.example.claude-usage.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.claude-usage</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>__HELPER_PATH__</string>
    </array>

    <!-- Run on load and then every 5 minutes. The /oauth/usage numbers move
         slowly, so this is plenty fresh and stays well within rate limits. -->
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>300</integer>

    <key>StandardOutPath</key>
    <string>__LOG_DIR__/helper.log</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/helper.err.log</string>
</dict>
</plist>
